# supervisord.conf
[supervisord]
nodaemon=true

[program:unbound]
command=/usr/bin/svscan /service
autorestart=true
